# scripts/supervisord-lamp.conf
[program:mysqld]
command=/scripts/mysql-start.sh
numprocs=1
autostart=true
autorestart=true

[program:apache2]
command=/scripts/apache2-start.sh
numprocs=1
autostart=true
autorestart=true

